# skills/recruit-ops/deploy/systemd/recruit-cron-runner.timer
[Unit]
Description=Recruit Ops cron orchestrator — every 10 minutes
Documentation=file:///home/admin/recruit-workspace/skills/recruit-ops/docs/CLI_REFERENCE.md

[Timer]
# v3.5 文档建议节奏：每 10 分钟一轮
# 每次 cron_runner 自己根据 only_hours 决定要不要跑 health_check
OnCalendar=*:0/10
Persistent=true
# 单元被启动后再随机 0-30s 抖动，避免和其他用户级 timer 抢 LLM/IMAP
RandomizedDelaySec=30
AccuracySec=30s

[Install]
WantedBy=timers.target
