# skills/recruit-ops/deploy/systemd/recruit-cron-runner.service
[Unit]
Description=Recruit Ops cron orchestrator (v3.5.10) — runs cron_runner once per tick
# 串 inbox.cmd_scan + inbox.cmd_analyze + common.cmd_interview_reminder
# + auto_reject.cmd_scan_exam_timeout + ops.cmd_health_check（仅 09:xx）
# 任务表定义：scripts/cron/cron_runner.py::_TASKS
Documentation=file:///home/admin/recruit-workspace/skills/recruit-ops/docs/CLI_REFERENCE.md

[Service]
Type=oneshot
WorkingDirectory=/home/admin/recruit-workspace/skills/recruit-ops
Environment=RECRUIT_WORKSPACE_ROOT=/home/admin/recruit-workspace
Environment=RECRUIT_CONFIG_DIR=/home/admin/recruit-workspace/config
Environment=RECRUIT_DATA_ROOT=/home/admin/recruit-files
Environment=RECRUIT_EXAM_ARCHIVE_DIR=/home/admin/recruit-files/exam_package
Environment=PYTHONPATH=/home/admin/recruit-workspace/skills/recruit-ops/scripts
ExecStart=/home/admin/recruit-workspace/skills/recruit-ops/.venv/bin/python -m cron.cron_runner
# 5 个子任务上限 5×300s=25min；给 systemd 留 30 分钟兜底，超时杀进程
TimeoutStartSec=1800
# cron_runner 退出码：
#   0 = 全部子任务成功
#   1 = 至少一个子任务失败，但 _alert_boss 已经推飞书告警了（"软失败"）
# 把 1 也算 success，避免 systemctl status 长期挂 failed 干扰真正排障；
# 真 Python crash / ImportError / OOM 会拿到非 0/1 的退出码或被信号杀，
# systemd 仍会标 failed。失败信号的真源头是飞书 [CRON FAIL] 卡片，不是 systemctl。
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
